# stdlib/files/node_exporter.conf
#Prometheus Node Exporter Upstart script
start on startup
script
  /opt/node_exporter/node_exporter
end script
