# stdlib/files/prometheus.conf
#Prometheus Upstart Script
start on startup
script
  /opt/prometheus/prometheus -config.file /etc/prometheus.yml > /var/log/prometheus.log 2>&1 &
end script
