# files/alertmanager.conf
#Prometheus Node Exporter Upstart script
start on startup
script
  /opt/alertmanager-0.13.0.linux-amd64/alertmanager --config.file /opt/alertmanager-0.13.0.linux-amd64/alertmanager.yml 
end script
