# files/cloudwatch-exporter.conf
start on startup

script
  /usr/bin/java -jar /usr/local/bin/cloudwatch-exporter.jar 9106 /etc/cloudwatch-metrics.yml
end script
