# files/node_exporter.conf
#Prometheus Node Exporter Upstart script
start on startup
script
  /opt/node_exporter-0.15.2.linux-amd64/node_exporter
end script
