# files/prometheus.conf
#Prometheus Upstart Script
start on startup
script
  /opt/prometheus-2.0.0.linux-amd64/prometheus --config.file /opt/prometheus-2.0.0.linux-amd64/prometheus.yml
end script
